# Semestr_VI/Programowanie_Uslug_Sieciowych/PUS-08-SystemD/systemd/pus_sock.socket
[Unit]
Description=PUS Socket
PartOf=pus_sock.service

[Socket]
ListenStream=3301

[Install]
WantedBy=sockets.target
